# server/privacy-guardrail.service.example
# Example systemd unit for the FastAPI app.
# You cannot edit the real unit yourself (/etc/systemd/system/privacy-guardrail.service) —
# send these values to Rainer so he can adjust the existing service.
#
# The two things that must change vs. the dummy app:
#   - WorkingDirectory -> the folder that contains main.py AND dist/
#   - ExecStart        -> run THIS app (main:app) on localhost:50000
#
# Assumes a virtualenv at ~/web/privacy-guardrail/.venv (see deploy guide).

[Unit]
Description=Privacy Guardrail data-donation app
After=network.target

[Service]
User=privacy-guardrail
WorkingDirectory=/home/privacy-guardrail/web/privacy-guardrail
Environment=SUBMISSION_DIR=/home/privacy-guardrail/web/privacy-guardrail/data/submissions
ExecStart=/home/privacy-guardrail/web/privacy-guardrail/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 50000
Restart=on-failure

[Install]
WantedBy=multi-user.target
